# _etc/supervisord.conf
[program:ca_mylesb_draft]
directory = /home/myles/WebSites/ca_mylesb_www
command = fab build_draft
autostart = true
autorestart = true